# supervisor.conf
[supervisord]

[program:mysql]
command=/usr/bin/mysqld_safe
    --bind-address=0.0.0.0
stopsignal=KILL
[program:php]
command=/usr/bin/php5
    -S 0.0.0.0:80
directory=/netrunnerdb
